# controller/systemctl/arti-failure-notify.service
[Unit]
Description=Arti Failure Notification
After=arti.service

[Service]
Type=oneshot
ExecStart=/home/ec2-user/go/src/fcbh-dataset-io/controller/systemctl/arti-failure-notify.sh

StandardOutput=journal
StandardError=journal
SyslogIdentifier=arti-failure-notify